# grav-webdav/provision/templates/etc/supervisor/conf.d/cron.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/cron -f
startsecs=1
autostart=true
autorestart=unexpected
exitcodes=0
