# wordpress/provision/templates/etc/supervisor/conf.d/wp-cron.conf.j2
[program:wp-cron]
command=/usr/bin/schedule_program --minutes 15 --wake-up 30 -- /usr/bin/php /container/www/wp-cron.php
autostart=true
autorestart=unexpected
directory=/container/www/
