# gogs/provision/templates/etc/supervisor/conf.d/sshd.conf
[supervisord]
nodaemon=true

[program:sshd]
command=bash -c "/usr/sbin/sshd -D -f /etc/ssh/sshd_config"
startsecs=2
autostart=true
autorestart=unexpected
exitcodes=0
