# gogs/provision/templates/etc/supervisor/conf.d/sshd.conf.j2
[program:sshd]
command=/usr/sbin/sshd -D -f /etc/ssh/sshd_config
startsecs=2
autostart=true
autorestart=unexpected
exitcodes=0
redirect_stderr=true
