# mariadb-src/provision/templates/etc/supervisor/conf.d/mariadb.conf
[supervisord]
nodaemon=true

[program:mariadb]
command=/usr/sbin/mysqld
autostart=true
autorestart=false
startsecs=10
