# redis-src/provision/templates/etc/supervisor/conf.d/redis.conf
[supervisord]
nodaemon=true

[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
directory=/var/lib/redis
user=redis
autostart=true
autorestart=false
startsecs=2
